# dist/init/linux-upstart/caddy.conf
description "Caddy startup script"
author "Mathias Beke"

start on runlevel [2345]
stop on runlevel [016]


setuid www-data
setgid www-data

respawn
respawn limit 10 5

limit nofile 4096 4096

script
	exec /usr/bin/caddy -agree=true -conf=/etc/caddy/Caddyfile
end script